# src/bin/raspibot.service
[Unit]
Description=RaspiBOT for Telegram

[Service]
User=SYS_USER
Type=oneshot
ExecStart=/bin/raspibot exec
TimeoutSec=0
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target